# templates/tempest_run.service.j2
[Service]
Type=oneshot
User={{ rally_user }}
ExecStart=/usr/bin/sh -c '{{ deployments_folder }}/{{ item.key }}-tempest_run.sh --tag cron_$(date +%Y_%m_%d_%H)'
