# templates/tempest_run.timer.j2
[Timer]
OnCalendar=*-*-* {{ item.value.tempest_cron_run_hour }}:00:00
Unit=tempest_run_{{ item.key }}.service

[Install]
WantedBy=timers.target
